# files/supervisord-sensu-api.conf
[program:sensu-api]
command=/opt/sensu/bin/sensu-api -c /etc/sensu/settings.json -d /etc/sensu/conf.d
